# mkosi.extra/usr/lib/systemd/system/nix-daemon.socket.d/10-nix-overlay.conf
[Unit]
After=nix-overlay.service
Requires=nix-overlay.service
