# subprojects/bluefin-common/system_files/shared/etc/profile.d/ublue-fastfetch.sh
alias neofetch='ublue-fastfetch'
alias neowofetch='ublue-fastfetch'
alias fastfetch='ublue-fastfetch'
