# subprojects/bluefin-common/system_files/shared/etc/profile.d/ublue-motd.sh
#!/usr/bin/env bash

[ ! -e "$HOME"/.config/no-show-user-motd ] && ublue-motd
