# mkosi.extra/usr/lib/systemd/system/miasma-set-hostname.service
[Unit]
Description=Set default Miasma hostname on first boot
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ConditionPathExists=!/var/lib/miasma/.hostname-initialized
ExecStart=/bin/bash -c 'echo miasma > /etc/hostname && hostname miasma'
ExecStartPost=/bin/bash -c 'mkdir -p /var/lib/miasma && touch /var/lib/miasma/.hostname-initialized'

[Install]
WantedBy=sysinit.target
