# mkosi.extra/etc/profile.d/mfetch.sh
alias neofetch='/usr/bin/mfetch'
alias neowofetch='/usr/bin/mfetch'
alias fastfetch='/usr/bin/mfetch'
